# deploy/supervisor.conf
; ─── HotelHub HMS — Queue Worker ──────────────────────────────

[program:hotel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/hotel/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --timeout=300
autostart=true
autorestart=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/hotel-queue.log
stderr_logfile=/var/log/supervisor/hotel-queue-error.log
stopwaitsecs=180

; ─── HotelHub HMS — Scheduler ──────────────────────────────────

[program:hotel-scheduler]
command=php /var/www/hotel/artisan schedule:work
user=www-data
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/hotel-scheduler.log
stderr_logfile=/var/log/supervisor/hotel-scheduler-error.log
stopwaitsecs=60
