# scripts/license-server/supervisord.conf
; HotelHub License Server — Supervisor Configuration
; Manages: PHP-FPM, Nginx, Queue Workers (3 processes), Scheduler

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/run/supervisord.pid
user=root

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ─────────────────────────────────────────────────────────────────
; Program: PHP-FPM (priority 5)
; ─────────────────────────────────────────────────────────────────
[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
priority=5
startretries=5
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─────────────────────────────────────────────────────────────────
; Program: Nginx (priority 10)
; ─────────────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
startretries=5
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─────────────────────────────────────────────────────────────────
; Program: Queue Worker — Default queue (priority 20)
; ─────────────────────────────────────────────────────────────────
[program:queue-default]
command=php /var/www/artisan queue:work --queue=default --sleep=3 --tries=3 --max-time=3600 --max-jobs=500
autostart=true
autorestart=true
user=www-data
priority=20
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=60
stopsignal=SIGTERM

; ─────────────────────────────────────────────────────────────────
; Program: Queue Worker — Notifications (priority 21)
; ─────────────────────────────────────────────────────────────────
[program:queue-notifications]
command=php /var/www/artisan queue:work --queue=notifications --sleep=2 --tries=3 --max-time=1800
autostart=true
autorestart=true
user=www-data
priority=21
numprocs=1
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=30
stopsignal=SIGTERM

; ─────────────────────────────────────────────────────────────────
; Program: Queue Worker — Heartbeat (priority 22)
; ─────────────────────────────────────────────────────────────────
[program:queue-heartbeat]
command=php /var/www/artisan queue:work --queue=heartbeat --sleep=1 --tries=1 --max-time=600
autostart=true
autorestart=true
user=www-data
priority=22
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=30
stopsignal=SIGTERM

; ─────────────────────────────────────────────────────────────────
; Program: Laravel Scheduler (priority 30)
; Runs artisan schedule:run every 60 seconds
; ─────────────────────────────────────────────────────────────────
[program:scheduler]
command=/bin/sh -c "while true; do php /var/www/artisan schedule:run --no-interaction >> /dev/null 2>&1; sleep 60; done"
autostart=true
autorestart=true
user=www-data
priority=30
numprocs=1
stdout_logfile=/var/log/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=10
stopsignal=SIGTERM
